# roles/mountshare/templates/etc/systemd/system/mountshare.automount.j2
[Unit]
Description=Automount mountshare

[Automount]
Where={{ mountshare_mountpoint }}

[Install]
WantedBy=multi-user.target
